# RuneAudio+RuneUIe.img/runonce.service
[Unit]
After=udevil.service

[Service]
Type=idle
ExecStartPre=/bin/sleep 5
ExecStart=/root/runonce.sh

[Install]
WantedBy=multi-user.target
